# launchd/com.boyoclaw.runtime.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  BoyoClaw / MyClaw — LaunchAgent (user session)

  1. Copy to ~/Library/LaunchAgents/com.boyoclaw.runtime.plist
  2. Edit REPO_ROOT and PYTHON_BIN to your machine (absolute paths only).
  3. launchctl load ~/Library/LaunchAgents/com.boyoclaw.runtime.plist
  4. launchctl start com.boyoclaw.runtime

  Stop:  launchctl stop com.boyoclaw.runtime
  Unload: launchctl unload ~/Library/LaunchAgents/com.boyoclaw.runtime.plist

  Requires: Ollama (or your model stack) available; telegram.json in .sandbox/workspace/
  Use --telegram --headless so launchd does not attach a TTY (stdin EOF would exit otherwise).
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.boyoclaw.runtime</string>

	<key>WorkingDirectory</key>
	<string>REPO_ROOT</string>

	<key>ProgramArguments</key>
	<array>
		<string>PYTHON_BIN</string>
		<string>-m</string>
		<string>runtime</string>
		<string>--telegram</string>
		<string>--headless</string>
	</array>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>src</string>
		<!-- Optional: if Ollama is not default -->
		<!-- <key>OLLAMA_HOST</key> -->
		<!-- <string>http://127.0.0.1:11434</string> -->
	</dict>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/Users/YOU/Library/Logs/BoyoClaw/runtime.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOU/Library/Logs/BoyoClaw/runtime.err.log</string>
</dict>
</plist>
